# configuration/scheduler/supervisor-worker.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.out.log

[program:php]
command=php /app/your-php-script.php
autostart=true
autorestart=true
stderr_logfile=/var/log/php.err.log
stdout_logfile=/var/log/php.out.log
